# files/etc/uci-defaults/99-set-kucat-theme
#!/bin/sh
# Set kucat as default LuCI theme
uci set luci.main.mediaurlbase='/luci-static/kucat'
uci commit luci
